# systemd/rc-bot.service
[Unit]
Description=Rocket.Chat Infrastructure Assistant Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/data/local/infra-assistant/rc-bot
EnvironmentFile=/data/local/infra-assistant/.env
ExecStart=/data/local/infra-assistant/rc-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=append:/data/local/infra-assistant/logs/rc-bot.log
StandardError=append:/data/local/infra-assistant/logs/rc-bot.log

[Install]
WantedBy=multi-user.target
