# systemd/zabbix-poller.service
[Unit]
Description=Zabbix Alert Poller (one-shot)
After=network-online.target

[Service]
Type=oneshot
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/data/local/infra-assistant/zabbix-poller
EnvironmentFile=/data/local/infra-assistant/.env
ExecStart=/data/local/infra-assistant/zabbix-poller/venv/bin/python poller.py
StandardOutput=append:/data/local/infra-assistant/logs/zabbix-poller.log
StandardError=append:/data/local/infra-assistant/logs/zabbix-poller.log
