# systemd/zammad-poller.service
[Unit]
Description=Zammad Unassigned Ticket Poller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/data/local/infra-assistant/zammad-poller
EnvironmentFile=/data/local/infra-assistant/.env
ExecStart=/data/local/infra-assistant/zammad-poller/venv/bin/python poller.py
Restart=always
RestartSec=10
StandardOutput=append:/data/local/infra-assistant/logs/zammad-poller.log
StandardError=append:/data/local/infra-assistant/logs/zammad-poller.log

[Install]
WantedBy=multi-user.target
